# package/a2a-gateway-test/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        export A2A_SKIP_GIT=1
        export A2A_SKIP_APT=1

        # Override settings for the test instance
        export A2A_DOMAIN="test.greenhelix.net"
        export A2A_WWW_DOMAIN=""
        export PORT=8001
        export INSTALL_DIR="/opt/a2a-test"
        export DATA_DIR="/var/lib/a2a-test"

        # Create test data directories
        mkdir -p /var/lib/a2a-test
        mkdir -p /var/log/a2a-test

        # Create a2a system user if not exists
        if ! id -u a2a >/dev/null 2>&1; then
            useradd --system --shell /usr/sbin/nologin --home-dir /opt/a2a-test a2a
        fi

        chown -R a2a:a2a /var/lib/a2a-test /var/log/a2a-test /opt/a2a-test

        # Run database migrations
        if [[ -x /opt/a2a-test/scripts/migrate_db.sh ]]; then
            A2A_DATA_DIR=/var/lib/a2a-test /opt/a2a-test/scripts/migrate_db.sh || exit 1
        fi

        # Set up Python venv and install deps
        python3.12 -m venv /opt/a2a-test/venv
        /opt/a2a-test/venv/bin/pip install --upgrade pip wheel setuptools -q
        /opt/a2a-test/venv/bin/pip install \
            'fastapi>=0.115' 'uvicorn>=0.29' 'httpx>=0.27' \
            'aiosqlite>=0.20' 'pydantic>=2.0' 'cryptography>=46.0' \
            'z3-solver>=4.12' -q
        /opt/a2a-test/venv/bin/pip install -e /opt/a2a-test/sdk -q 2>/dev/null \
            || /opt/a2a-test/venv/bin/pip install /opt/a2a-test/sdk -q

        # Verify z3-solver is functional (C extension must compile)
        /opt/a2a-test/venv/bin/python -c "import z3; print('z3-solver OK:', z3.get_version_string())"

        # Create .env if missing
        if [[ ! -f /opt/a2a-test/.env ]]; then
            cat > /opt/a2a-test/.env << 'ENVEOF'
A2A_ENV=staging
HOST=127.0.0.1
PORT=8001
A2A_DATA_DIR=/var/lib/a2a-test
BILLING_DSN=sqlite:////var/lib/a2a-test/billing.db
PAYWALL_DSN=sqlite:////var/lib/a2a-test/paywall.db
PAYMENTS_DSN=sqlite:////var/lib/a2a-test/payments.db
MARKETPLACE_DSN=sqlite:////var/lib/a2a-test/marketplace.db
TRUST_DSN=sqlite:////var/lib/a2a-test/trust.db
EVENT_BUS_DSN=sqlite:////var/lib/a2a-test/events.db
WEBHOOK_DSN=sqlite:////var/lib/a2a-test/webhooks.db
IDENTITY_DSN=sqlite:////var/lib/a2a-test/identity.db
MESSAGING_DSN=sqlite:////var/lib/a2a-test/messaging.db
DISPUTE_DSN=sqlite:////var/lib/a2a-test/disputes.db
LOG_LEVEL=INFO
FORCE_HTTPS=1
# --- Formal Gatekeeper (in-process Z3) ---
VERIFIER_AUTH_MODE=mock
ENVEOF
            chmod 600 /opt/a2a-test/.env
        fi

        # Ensure critical env vars exist in .env (fixes audit H2/C1 for upgrades)
        _env_file="/opt/a2a-test/.env"
        if [[ -f "$_env_file" ]]; then
            grep -q '^FORCE_HTTPS=' "$_env_file" || echo 'FORCE_HTTPS=1' >> "$_env_file"
            grep -q '^A2A_ENV=' "$_env_file" || echo 'A2A_ENV=staging' >> "$_env_file"
            if grep -q '^VERIFIER_AUTH_MODE=' "$_env_file"; then
                sed -i 's/^VERIFIER_AUTH_MODE=.*/VERIFIER_AUTH_MODE=mock/' "$_env_file"
            else
                echo 'VERIFIER_AUTH_MODE=mock' >> "$_env_file"
            fi
        fi

        # --- Cloudflare Authenticated Origin Pulls ---
        CF_CA_PATH="/etc/nginx/cloudflare-origin-pull-ca.pem"
        CF_CA_URL="https://developers.cloudflare.com/ssl/static/authenticated_origin_pull_ca.pem"
        if [[ ! -f "$CF_CA_PATH" ]]; then
            echo "Downloading Cloudflare Authenticated Origin Pull CA..."
            curl -fsSL "$CF_CA_URL" -o "$CF_CA_PATH" \
                || echo "WARNING: Failed to download Cloudflare origin-pull CA."
            [[ -f "$CF_CA_PATH" ]] && chmod 644 "$CF_CA_PATH"
        fi

        # Enable the a2a-test site in nginx
        ln -sf /etc/nginx/sites-available/a2a-test /etc/nginx/sites-enabled/a2a-test
        nginx -t && systemctl reload nginx

        # systemd service
        systemctl daemon-reload
        systemctl enable a2a-gateway-test
        systemctl restart a2a-gateway-test
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
