# package/a2a-db-backup/etc/systemd/system/a2a-db-backup.timer
[Unit]
Description=Run A2A database backup hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
