# package/a2a-gateway/etc/systemd/system/a2a-gateway.service
[Unit]
Description=A2A Commerce Gateway
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=a2a
Group=a2a
WorkingDirectory=/opt/a2a
EnvironmentFile=/opt/a2a/.env
ExecStart=/opt/a2a/venv/bin/python -m uvicorn gateway.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --log-level info \
    --access-log
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# --- Secrets management (systemd credentials) ---
# Uncomment to migrate secrets from .env to per-service credential files.
# Create credential files in /etc/credstore/a2a/ (mode 0600, owner root:root).
# The service reads them via $CREDENTIALS_DIRECTORY/<name> at runtime.
# See docs/policies/secrets-management-policy.md for migration steps.
#LoadCredential=stripe_secret_key:/etc/credstore/a2a/stripe_secret_key
#LoadCredential=stripe_webhook_secret:/etc/credstore/a2a/stripe_webhook_secret
#LoadCredential=jwt_secret:/etc/credstore/a2a/jwt_secret
#LoadCredential=admin_key:/etc/credstore/a2a/admin_key

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/a2a /var/log/a2a
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
