# package/a2a-gateway/etc/systemd/system/a2a-gateway.service
[Unit]
Description=A2A Commerce Gateway
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=a2a
Group=a2a
WorkingDirectory=/opt/a2a
EnvironmentFile=/opt/a2a/.env
ExecStart=/opt/a2a/venv/bin/python -m uvicorn gateway.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --log-level info \
    --access-log
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/a2a /var/log/a2a
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
